# package/skylark_daemon/overlay/etc/sv/skylark_http2_runsvdir/run
#!/bin/sh

name=skylark_http2
service=${name}_runsvdir
user=skylarkd

source /etc/init.d/template_svdir.sh
