# package/cell_modem_daemon/overlay/etc/sv/cell_modem_daemon_runsvdir/run
#!/bin/ash

name=cell_modem_daemon
service=${name}_runsvdir
user=celld

source /etc/init.d/template_svdir.sh
