# package/ota_daemon/overlay/etc/sv/ota_daemon_runsvdir/run
#!/bin/sh

name=ota_daemon
service=${name}_runsvdir
user=otad

source /etc/init.d/template_svdir.sh
