# board/piksiv3/rootfs-overlay/usr/lib/systemd/system/load-fpga-image.service
[Unit]
Description=Load FPGA image
Wants=dev-xdevcfg.device
After=dev-xdevcfg.device

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/swiftnav/load_fpga_image
