# board/piksiv3/rootfs-overlay/usr/lib/systemd/system/persistent.mount
[Unit]
Description=Mount persistent partition

[Mount]
What=/dev/mtdblock_persistent
Where=/persistent
Type=jffs2
Options=defaults
