# board/piksiv3/rootfs-overlay/usr/lib/systemd/system/piksi-leds.service
[Unit]
Description=Piksi LEDs

[Service]
ExecStart=/bin/piksi_leds

[Install]
WantedBy=multi-user.target
