# board/piksiv3/rootfs-overlay/usr/lib/systemd/system/piksi-usb-gadget.service
[Unit]
Description=Configure usb gadget

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/swiftnav/piksi-usb-gadget
