# board/piksiv3/rootfs-overlay/usr/lib/systemd/system/zmq_router@.service
[Unit]
Description=zmq router

[Service]
ExecStart=/bin/zmq_router -f /etc/zmq_router/%i.yml

[Install]
WantedBy=multi-user.target
