# package/ports_daemon/overlay/etc/sv/ports_daemon_runsvdir/run
#!/bin/sh
 exec 2>&1
 exec chpst -u portsd runsvdir /var/run/ports_daemon/sv
