# package/ports_daemon/overlay/etc/sv/ports_daemon_runsvdir/run
#!/bin/sh

name=ports_daemon
service=${name}_runsvdir
user=portsd

source /etc/init.d/template_svdir.sh
