# package/resource_monitor/resource_monitor/overlay/etc/sv/resource_monitor_runsvdir/run
#!/bin/sh

name=resource_monitor
user=resmond
service=${name}_runsvdir

source /etc/init.d/template_svdir.sh
