# package/piksi_system_daemon/overlay/etc/sv/piksi_system_daemon_runsvdir/run
#!/bin/sh

name=piksi_system_daemon
user=piksi_sys
service=${name}_runsvdir

source /etc/init.d/template_svdir.sh
