# package/skylark_daemon/overlay/etc/sv/skylark_adapt_ul_runsvdir/run
#!/bin/sh

name=skylark_adapt_ul
service=${name}_runsvdir
user=skylarkd

source /etc/init.d/template_svdir.sh
